# auto-email-worker.service
[Unit]
Description=Auto Email Sender - RabbitMQ Queue Worker
After=network.target

[Service]
Type=simple
User=payatu
WorkingDirectory=/home/payatu/personal/auto_email_sender
ExecStart=/usr/bin/python3 /home/payatu/personal/auto_email_sender/worker.py
Restart=always
RestartSec=10
StandardOutput=append:/home/payatu/personal/auto_email_sender/worker.log
StandardError=append:/home/payatu/personal/auto_email_sender/worker.log

[Install]
WantedBy=multi-user.target
